# scripts/agentdesk-browser-tools-server.service
[Unit]
Description=AgentDesk Browser Tools Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=%h/.asdf/shims/npx @agentdeskai/browser-tools-server@1.2.0
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
